# supervisord_lightweight.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/supervisord.pid
user=root

[program:fastapi]
command=python /app/fastapi_server_lightweight.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/fastapi.err.log
stdout_logfile=/app/logs/fastapi.out.log
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
priority=1

[program:signal_worker]
command=python /app/signal_worker_lightweight.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/signal_worker.err.log
stdout_logfile=/app/logs/signal_worker.out.log
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",SIGNAL_SCAN_INTERVAL="60"
priority=2
startsecs=10